# ircd-hybrid/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/ircd-hybrid
user=irc

